# dists/jetty-deb/src/main/unix/scripts/prerm
#!/bin/sh
set -e

if [ -x "/etc/init.d/jetty" ]; then
    invoke-rc.d jetty stop || exit $?
fi
